# deploy/digitalocean/files/etc/systemd/system/swe-swe.service
[Unit]
Description=swe-swe AI Development Environment
After=cloud-init.target docker.service network-online.target
Requires=docker.service
Wants=cloud-init.target network-online.target

[Service]
Type=simple
User=swe-swe
Group=swe-swe
EnvironmentFile=/etc/swe-swe/env
WorkingDirectory=/workspace
ExecStart=/usr/local/bin/swe-swe up --project-directory /workspace
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
